# systemd/piper-proxy.path
[Unit]
Description=Monitor inbox for new Piper jobs

[Path]
DirectoryNotEmpty=%h/Code/Piper/inbox
Unit=piper-proxy.service

[Install]
WantedBy=paths.target
